# service/clamshell.openrc
#!/sbin/openrc-run

description="Clamshell mode suspend daemon"

command="/usr/bin/clamshell"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need localmount
  after display-manager
}
